# systemd/ntpclock.service
[Unit]
Description=One-shot NTP clock sync (ntpclock)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ntpclock
StandardOutput=journal
StandardError=journal
